# config/init/nginx.conf
# nginx

stop on runlevel [016]
start on runlevel [345]

respawn
instance /usr/local/nginx/sbin/nginx
exec /etc/init.d/nginx start
