# config/init/repos.conf
# git-http-backend.go

stop on runlevel [016]
start on runlevel [345]

respawn
instance /app/gitcache.org/bin/repos
exec su gitcache -c "/app/gitcache.org/bin/repos >> /tmp/repos.log 2>&1"
